# sv/update-label/run
#!/bin/bash

source /etc/envvars
env

IP=`ip route get 1 | awk '{print $NF;exit}'`
HOST=`hostname -s`

exec 2>&1
while true; do kubectl label nodes ${IP} --overwrite host=${HOST}; sleep 60; done;
