# sv/register/run
#!/bin/bash

source /etc/envvars
env

exec 2>&1
while true; do /register.sh; sleep 30; done;
