# sv/register/run
#!/bin/bash

source /etc/envvars
env

/register.sh

exec 2>&1
exec sleep infinity
