# sv/confd/run
#!/bin/bash

#sv start register  || exit 1

source /etc/envvars
env

exec 2>&1
exec confd -watch=true -interval 30 -backend=etcd -node "http://$ETCD_HOST:4001" -log-level="info"

