# sv/kubelet/run
#!/bin/bash

source /etc/envvars
env

KUBELET_OPTS="--address=0.0.0.0 \
--port=10250 \
--etcd_servers=http://${ETCD_HOST}:4001 \
--logtostderr=true"

exec 2>&1
exec kubelet ${KUBELET_OPTS}
